# home/crontab
17 9 * * 1 PATH="$HOME/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin" "$HOME/bin/,ai-kb" curate >> "$HOME/.local/share/ai-kb/curate.log" 2>&1
7 9 * * * /usr/bin/rsync -a /tmp/specs/ "$HOME/.local/share/agent-specs-archive/" >/dev/null 2>&1
